# deploy/systemd/analytics-sync-worker.service
[Unit]
Description=amoCRM Analytics Sync Worker
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
WorkingDirectory=/opt/analytics/apps/api
Environment=NODE_ENV=production
Environment=NODE_OPTIONS=--max-old-space-size=512
ExecStart=/usr/bin/node dist/sync-worker.js
Restart=always
RestartSec=5
User=ameri
Group=ameri
MemoryHigh=650M
MemoryMax=800M
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
